# systemd/genomeinsight-api.service
[Unit]
Description=GenomeInsight API Server (Uvicorn)
After=network.target

[Service]
Type=simple
ExecStart=uvicorn backend.main:app --host 127.0.0.1 --port 8000
WorkingDirectory=%h/GenomeInsight
Restart=on-failure
Environment=PATH=%h/.local/bin:/usr/bin

[Install]
WantedBy=default.target
